# examples/python-app/supervisord.conf
; Four worker replicas, each binding its OWN unix socket under /sockets (a
; volume shared with the edge-peer container). %(process_num)02d makes every
; replica's socket file distinct — no SO_REUSEPORT-on-a-file problem.
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:app]
command=python3 /app/%(ENV_APP)s /sockets/app_%(process_num)02d.sock W%(process_num)02d 1 %(ENV_MUX_TOKEN)s
process_name=%(program_name)s_%(process_num)02d
numprocs=4
autorestart=true
startsecs=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
