# deploy/supervisord.conf
; Run N worker replicas, each binding its OWN unix socket, with supervisord.
; %(process_num)02d gives every replica a distinct socket file, so there is no
; SO_REUSEPORT-on-a-file problem — the workers are servers and the edge-peer
; dials each one. List the same socket files in the edge-peer command (below).
;
;   supervisord -c deploy/supervisord.conf
;
; Env: MUX_TOKEN is the pre-shared key shared by edge, edge-peer and workers.

[supervisord]
nodaemon=true
logfile=/tmp/mux-supervisord.log

[program:worker]
command=python3 %(here)s/../python/worker.py /run/mux/app_%(process_num)02d.sock W%(process_num)02d 4 %(ENV_MUX_TOKEN)s
process_name=%(program_name)s_%(process_num)02d
numprocs=4
directory=%(here)s/../python
autorestart=true
stdout_logfile=/tmp/mux-worker-%(process_num)02d.log
stderr_logfile=/tmp/mux-worker-%(process_num)02d.err
startsecs=1

; The edge-peer dials the edge upstream and every worker socket. Because each
; replica has a fixed name, the worker list is static and explicit here:
[program:edge-peer]
command=edge-peer --connect %(ENV_EDGE_ADDR)s
    --worker /run/mux/app_00.sock
    --worker /run/mux/app_01.sock
    --worker /run/mux/app_02.sock
    --worker /run/mux/app_03.sock
    --token %(ENV_MUX_TOKEN)s
autorestart=true
stdout_logfile=/tmp/mux-edge-peer.log
stderr_logfile=/tmp/mux-edge-peer.err
startsecs=1
